# scripts/com.codeflicker.ai-weekly-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  AI周报定时推送任务
  触发时间: 每周一 10:00
  执行链路: launchd → 个人助理_V1/scripts/run-ai-weekly-report.sh (委托入口)
           → AI-Insight/scripts/run-ai-weekly-report.sh → send_ai_weekly.py
  注意: 实际加载的plist在 ~/Library/LaunchAgents/ 中，由个人助理_V1统一管理。
        本文件仅作为参考文档，与实际加载版本保持一致。
  安装: cp scripts/com.cf.ai-weekly-report.plist ~/Library/LaunchAgents/
        launchctl load ~/Library/LaunchAgents/com.cf.ai-weekly-report.plist
-->
<dict>
    <key>Label</key>
    <string>com.cf.ai-weekly-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shenlang/Documents/Codeflicker/个人助理_V1/scripts/run-ai-weekly-report.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <!-- 每周一 (Weekday=1) 上午10点 -->
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/shenlang/Documents/Codeflicker/个人助理_V1/data/logs/ai-weekly-report-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shenlang/Documents/Codeflicker/个人助理_V1/data/logs/ai-weekly-report-stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/shenlang/Documents/Codeflicker/个人助理_V1</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>LANG</key>
        <string>zh_CN.UTF-8</string>
        <key>HOME</key>
        <string>/Users/shenlang</string>
    </dict>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
