# scripts/com.codeflicker.ai-daily-report.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.codeflicker.ai-daily-report</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/shenlang/Documents/Codeflicker/AI-Insight/scripts/run-ai-daily-report.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/shenlang/Documents/Codeflicker/AI-Insight/scripts/logs/launchd-daily.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shenlang/Documents/Codeflicker/AI-Insight/scripts/logs/launchd-daily-error.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/shenlang/Documents/Codeflicker/AI-Insight</string>
</dict>
</plist>
